# nicepatrol-queue.service
[Unit]
Description=Nice Patrol Queue Worker
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
Restart=always
RestartSec=5s
ExecStart=/usr/bin/php /path/to/your/project/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --timeout=300
WorkingDirectory=/path/to/your/project
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target